# bin/cron.daily/andr-build
#!/usr/bin/env bash

su -c '
build_path="/home/andrew/andr-build"
[ -d "$build_path" ] || exit 1
pushd "$build_path"

for i in 1/1/{7,9}/{1,3}; do
	[ -d "$i" ] || continue
	pushd "$i"
	for j in *; do
		[ -d "$j" ] || continue
		pushd "$j"
		p=(`echo "$i" | sed "s/\// /g"`)
		andr-build "${p[@]}"
		popd
	done
	popd
done

popd
' - andrew

